# deploy/ccis-backend.service
[Unit]
Description=CCIS CodeHub Backend (Gunicorn)
After=network.target

[Service]
User=deploy
Group=www-data
WorkingDirectory=/home/deploy/CCIS-CodeHub/backend
Environment="PATH=/home/deploy/CCIS-CodeHub/backend/venv/bin"
EnvironmentFile=/home/deploy/CCIS-CodeHub/backend/.env
ExecStart=/home/deploy/CCIS-CodeHub/backend/venv/bin/gunicorn --workers 3 --bind 127.0.0.1:8000 core.wsgi:application
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
